# support/com.shazam.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shazam.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>__SHAZAM_CORE__/bin/start-daemon.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/.shazam/logs/daemon.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.shazam/logs/daemon.stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SHAZAM_DAEMON</key>
        <string>true</string>
        <key>SHAZAM_PORT</key>
        <string>4040</string>
    </dict>
</dict>
</plist>
